# misc/rc.local
#!/bin/bash

source /opt/ioi/config.sh
/opt/ioi/sbin/startatd.sh &

# vim: ft=sh ts=4 noet
