# misc/rc.local
#!/bin/bash

source /opt/icpcbo/config.sh
/opt/icpcbo/sbin/startatd.sh &

# vim: ft=sh ts=4 noet
